# .config/systemd/user/gtd_backup.service
[Unit]
Description="backup up important gtd files"

[Service]
ExecStart=/usr/bin/python /home/pipz/myconf/utils/bin/gtd_backup.py
